# plugins/appdata.backup/appdata_backup.cron
# Appdata.Backup cron settings
0 2 * * 0 php /usr/local/emhttp/plugins/appdata.backup/scripts/backup.php > /dev/null 2>&1
